# integration-tests/bot/tasks.ini
[supervisord]
pidfile = %(here)s/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file = %(here)s/supervisor.sock

[supervisorctl]
serverurl = unix://%(here)s/supervisor.sock

[group:node0]
programs = tx-enclave-node0,chain-abci-node0,tendermint-node0,client-rpc-node0

[program:tx-enclave-node0]
command = docker run --rm -p 26650:25933 --env RUST_BACKTRACE=1 --env RUST_LOG=info -v /home/leejw51/Github/mychain/integration-tests/bot/node0/enclave:/enclave-storage  integration-tests-chain-tx-enclave-sw
stdout_logfile = %(here)s/tx-enclave-%(group_name)s.log
environment = RUST_BACKTRACE=1,RUST_LOG=info
autostart = true
autorestart = true
redirect_stderr = true
priority = 0
startsecs = 1
startretries = 10

[program:chain-abci-node0]
command = chain-abci -g 3A76E0E2220DBA7D07742F29F625E70EC9FA79080A38C01A2A9AA70BBA305679 -c test-chain-y3m1e6-AB --enclave_server tcp://127.0.0.1:26650 --data /home/leejw51/Github/mychain/integration-tests/bot/node0/chain -p 26658
stdout_logfile = %(here)s/chain-abci-%(group_name)s.log
environment = RUST_BACKTRACE=1,RUST_LOG=info
autostart = true
autorestart = true
redirect_stderr = true
priority = 1
startsecs = 1
startretries = 10

[program:tendermint-node0]
command = tendermint node --home=/home/leejw51/Github/mychain/integration-tests/bot/node0/tendermint
stdout_logfile = %(here)s/tendermint-%(group_name)s.log
environment = RUST_BACKTRACE=1,RUST_LOG=info
autostart = true
autorestart = true
redirect_stderr = true
priority = 2
startsecs = 1
startretries = 10

[program:client-rpc-node0]
command = client-rpc --port=26651 --chain-id=test-chain-y3m1e6-AB --storage-dir=/home/leejw51/Github/mychain/integration-tests/bot/node0/wallet --websocket-url=ws://127.0.0.1:26657/websocket
stdout_logfile = %(here)s/client-rpc-%(group_name)s.log
environment = RUST_BACKTRACE=1,RUST_LOG=info
autostart = true
autorestart = true
redirect_stderr = true
priority = 3
startsecs = 1
startretries = 10

[group:node1]
programs = tx-enclave-node1,chain-abci-node1,tendermint-node1,client-rpc-node1

[program:tx-enclave-node1]
command = docker run --rm -p 26660:25933 --env RUST_BACKTRACE=1 --env RUST_LOG=info -v /home/leejw51/Github/mychain/integration-tests/bot/node1/enclave:/enclave-storage  integration-tests-chain-tx-enclave-sw
stdout_logfile = %(here)s/tx-enclave-%(group_name)s.log
environment = RUST_BACKTRACE=1,RUST_LOG=info
autostart = true
autorestart = true
redirect_stderr = true
priority = 0
startsecs = 1
startretries = 10

[program:chain-abci-node1]
command = chain-abci -g 3A76E0E2220DBA7D07742F29F625E70EC9FA79080A38C01A2A9AA70BBA305679 -c test-chain-y3m1e6-AB --enclave_server tcp://127.0.0.1:26660 --data /home/leejw51/Github/mychain/integration-tests/bot/node1/chain -p 26668
stdout_logfile = %(here)s/chain-abci-%(group_name)s.log
environment = RUST_BACKTRACE=1,RUST_LOG=info
autostart = true
autorestart = true
redirect_stderr = true
priority = 1
startsecs = 1
startretries = 10

[program:tendermint-node1]
command = tendermint node --home=/home/leejw51/Github/mychain/integration-tests/bot/node1/tendermint
stdout_logfile = %(here)s/tendermint-%(group_name)s.log
environment = RUST_BACKTRACE=1,RUST_LOG=info
autostart = true
autorestart = true
redirect_stderr = true
priority = 2
startsecs = 1
startretries = 10

[program:client-rpc-node1]
command = client-rpc --port=26661 --chain-id=test-chain-y3m1e6-AB --storage-dir=/home/leejw51/Github/mychain/integration-tests/bot/node1/wallet --websocket-url=ws://127.0.0.1:26667/websocket
stdout_logfile = %(here)s/client-rpc-%(group_name)s.log
environment = RUST_BACKTRACE=1,RUST_LOG=info
autostart = true
autorestart = true
redirect_stderr = true
priority = 3
startsecs = 1
startretries = 10

